# supervisor/systemd.service
[Unit]
Description=WordPress Sharp Image Processing Service
After=network.target
Wants=network.target

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart=/usr/local/bin/bun run start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
RestartPreventExitStatus=0
TimeoutSec=600
Environment=NODE_ENV=production
KillMode=mixed
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{UPLOADS_DIRECTORY}}
ReadWritePaths={{WORKING_DIRECTORY}}/logs

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target 